# static/usr/lib/udev/rules.d/90-fixrtc.rules
# when rtc* is added run fixrtc
ACTION=="add", SUBSYSTEM=="rtc", KERNEL=="rtc*", RUN+="/usr/lib/core/fixrtc"
