# static/usr/lib/udev/rules.d/90-fixrtc.rules
# when rtc* is added run fixrtc
SUBSYSTEM=="rtc", KERNEL=="rtc*", RUN+="/usr/bin/fixrtc"
