# static/usr/lib/systemd/system/secureboot-db.service.d/core-override.conf
[Unit]
ConditionKernelCommandLine=snapd_recovery_mode=install
