# static/usr/lib/systemd/system/var-lib-snapd-save.mount
[Unit]
ConditionPathIsMountPoint=/run/mnt/ubuntu-save

[Mount]
What=/run/mnt/ubuntu-save
Where=/var/lib/snapd/save
Options=bind
Type=none
