# static/usr/lib/systemd/system/systemd-bootchart.service.d/ubuntu-core.conf
[Unit]
Conflicts=shutdown.target
Before=shutdown.target

# poststop.sh needs to write to /run/mnt/data so we need to stop
# before read-only remounting happens
Conflicts=remount-data.service
Before=remount-data.service
# In case of install/recover/factory-reset, we want
# poststop.sh to stop before /run/mnt/ubuntu-data
After=run-mnt-ubuntu\x2ddata.mount

[Service]
KillSignal=SIGHUP
ExecStopPost=/lib/systemd/systemd-bootchart-poststop.sh
